# config/heisen_supervisor.conf
[program:gunicorn]
command=../venv/bin/gunicorn main:heisen -c heisen/wsgi.py
directory=/var/www/heisen_django/
user=nobody
autostart=true
autorestart=true
redirect_stderr=true
